# mu.timer
[Unit]
Description=Timer for mu indexing

[Timer]
OnUnitInactiveSec=1m

[Install]
WantedBy=mail.target
